# private_dot_config/systemd/user/screenshot.service
[Unit]
Description=Capture screen with maim

[Service]
Type=oneshot
# Adjust the command for your display, e.g., :0
Environment=DISPLAY=:0
ExecStartPre=/usr/bin/bash -c 'command -v maim >/dev/null 2>&1 || { echo "maim not found"; exit 1; }'
ExecStartPre=/usr/bin/mkdir -p /home/%u/screenshots
ExecStart=/usr/bin/bash -c '/usr/bin/maim "/home/%u/screenshots/screen_$(date +%%Y%%b%%d_%%H%%M%%S).png"'
