# private_dot_config/systemd/user/screenshot.timer
[Unit]
Description=Run screenshot capture every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
# Keep screenshots for 2 days (deletes files older than 2 days)
# This runs a find command once a day
Unit=screenshot.service

[Install]
WantedBy=timers.target
